# MWFA_Project/backend/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

################################################################
# البرنامج الأول: Kali-MCP Server (على المنفذ 8000 داخلياً)
# هذا هو "الدماغ" اللي يشغّل أدوات الأمان (nmap, etc.)
################################################################
[program:kali-mcp]
command=/usr/local/bin/kali-mcp-server
; إذا كان الأمر مختلف في الـ image جرّب:
; command=python3 /app/kali_mcp/server.py
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/kali-mcp.err.log
stdout_logfile=/var/log/supervisor/kali-mcp.out.log
priority=1
startsecs=5

################################################################
# البرنامج الثاني: Node.js Backend (على المنفذ 3001 للخارج)
# ينتظر 8 ثوانٍ حتى يبدأ Kali-MCP أولاً
################################################################
[program:node-backend]
command=sh -c "sleep 8 && npx prisma db push --skip-generate && node server.js"
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/node-backend.err.log
stdout_logfile=/var/log/supervisor/node-backend.out.log
environment=MCP_URL="http://localhost:8000",NODE_ENV="production"
priority=2
startsecs=3
